# roles/tracker/files/tracker.timer
[Unit]
Description=Run the tracker at 07,13:00:00

[Timer]
OnCalendar=*-*-* 07,13:00:00
Unit=tracker.service

[Install]
WantedBy=multi-user.target
