# roles/ims/templates/ims-backup-media.service.j2
[Unit]
Description=IMS daily media backup

[Service]
ExecStart=/var/projects/ebisc/scripts/backup_media.bash
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
