# roles/ims/templates/ims-updates.service.j2
[Unit]
Description=IMS daily imports
Requires=docker.service
Requires=ims-postgres.service
Requires=ims-elasticsearch.service
After=docker.service
After=ims-postgres.service
After=ims-elasticsearch.service

[Service]
Environment=HOSTNAME=%H
Environment=CONTAINER_NAME=%p
ExecStartPre=-/usr/bin/docker kill %p
ExecStartPre=-/usr/bin/docker kill mailer-%p
ExecStartPre=-/usr/bin/docker rm %p
ExecStartPre=-/usr/bin/docker rm mailer-%p
ExecStart=/usr/bin/docker run --name %p \
  --cap-drop=all \
  --user=$(id -u {{ ims_user }}) \
  -e DB_USER=${POSTGRESQL_USER} \
  -e DB_PASS=${POSTGRESQL_PASSWORD} \
  -e DB_HOST=ims-postgres \
  -e DB_NAME=${POSTGRESQL_DATABASE} \
  -e ES_HOST=ims-elasticsearch \
  -e SENDGRID_API_KEY=${SENDGRID_API_KEY} \
  -e SECRET_KEY=${DJANGO_SECRET_KEY} \
  -e HPSCREG_PASSWORD=${HPSCREG_PASSWORD} \
  -e HPSCREG_USER=${HPSCREG_USER} \
  -e LIMS_PASSWORD=${LIMS_PASSWORD} \
  -e LIMS_USER=${LIMS_USER} \
  -e BIOSAMPLE_KEY=${BIOSAMPLE_KEY} \
  -e SECRET_KEY=${DJANGO_SECRET_KEY} \
  --read-only \
  --tmpfs /tmp \
  --net=ebisc \
  -v {{ data_dir }}/media:/app/var/media:z \
  -v {{ data_dir }}/static:/app/var/static:z \
  ebisc/ims-uwsgi:latest \
  run-ims-update
ExecStop=/usr/bin/docker stop %p
TimeoutStartSec=0
ExecStopPost=/usr/bin/bash -c "/usr/bin/docker logs ${CONTAINER_NAME} 2>&1 \
  | /usr/bin/docker run -i \
    --name mailer-${CONTAINER_NAME} \
    --cap-drop=all \
    -e TO=${MAILTO} \
    -e SUBJECT=${CONTAINER_NAME} \
    -e FROM_HOST=${HOSTNAME} \
    -e AUTH_PASS=${SENDGRID_API_KEY} \
    ebisc/mailer:latest"

[Install]
WantedBy=multi-user.target
