# roles/ims/files/ims-updates.timer
[Unit]
Description=IMS updates at 04:00:00

[Timer]
OnCalendar=*-*-* 04:00:00
Unit=ims-updates.service

[Install]
WantedBy=multi-user.target
